# email_server/templates/dkimproxy_in.service.template
[Unit]
Description=DKIM Inbound Proxy
After=opensmtpd.service

[Service]
Type=simple
PIDFile=/usr/local/dkimproxy/var/run/dkimproxy_in.pid
ExecStart=/usr/local/dkimproxy/bin/dkimproxy.sh start-in
ExecStop=/usr/local/dkimproxy/bin/dkimproxy.sh stop-in

[Install]
WantedBy=multi-user.target
