# web/templates/myplaceonline-delayedjobs.service.template
[Unit]
Description=myplaceonline-delayedjobs
Requires=network-online.target
After=network-online.target

[Service]
WorkingDirectory=${cubevar_app_web_dir}

Environment="FTS_TARGET=${cubevar_app_full_text_search_target}" "RUBY_GC_MALLOC_LIMIT_MAX=${cubevar_app_rails_gc_max_newspace}" "RUBY_GC_OLDMALLOC_LIMIT_MAX=${cubevar_app_rails_gc_max_oldspace}" "RAILS_ENV=${cubevar_app_rails_environment}" "SECRET_KEY_BASE=${cubevar_app_passwords_devise_secret}" "ROOT_EMAIL=${cubevar_app_root_email}" "ROOT_PASSWORD=${cubevar_app_passwords_rails_root}" "SMTP_HOST=${cubevar_app_outbound_smtp}" "SMTP_USER=${cubevar_app_passwords_smtp_user}" "SMTP_PASSWORD=${cubevar_app_passwords_smtp_password}" "MAIL_FROM=${cubevar_app_mail_from}" "PERMDIR=${cubevar_app_nfs_client_mount}" "TWILIO_NUMBER=${cubevar_app_passwords_twilio_number}" "TWILIO_ACCOUNT=${cubevar_app_passwords_twilio_account}" "TWILIO_AUTH=${cubevar_app_passwords_twilio_auth}" "BUNDLE_GEMFILE=${cubevar_app_rails_gemfile}" "YELP_CONSUMER_KEY=${cubevar_app_passwords_yelp_consumer_key}" "YELP_CONSUMER_SECRET=${cubevar_app_passwords_yelp_consumer_secret}" "YELP_TOKEN=${cubevar_app_passwords_yelp_token}" "YELP_TOKEN_SECRET=${cubevar_app_passwords_yelp_token_secret}" "NODENAME=${cubevar_app_hostname}" "SOURCE_REVISION=${cubevar_app_source_revision}" "TRUSTED_CLIENTS=${cubevar_app_trusted_servers}" "GOOGLE_MAPS_API_KEY=${cubevar_app_passwords_google_maps_api_key}" "GOOGLE_MAPS_API_SERVER_KEY=${cubevar_app_passwords_google_maps_api_server_key}" "GOOGLE_PLACES_KEY=${cubevar_app_passwords_google_places_api_key}" "SECONDARY_GOOGLE_API_KEY=${cubevar_app_passwords_google_secondary_api_key}" "GEOLOCATE_KEY=${cubevar_app_passwords_geolocate_key}" "SPAWN_MAX_VSIZE=${cubevar_app_spawn_max_vsize}" "PAYPAL_CLIENT_ID=${cubevar_app_paypal_client_id}" "PAYPAL_SECRET=${cubevar_app_paypal_secret}" "MAILING_ADDRESS=${cubevar_app_mailing_address}" "BITCOIN_ADDRESS=${cubevar_app_bitcoin_address}" "MINCACHE=true" "BROWSER_PATH=${cubevar_app_chrome_path}" "LOADER_CACHE_DIR=/opt/cache/"

ExecStart=${cubevar_app_web_dir}/bin/delayed_job run -n 1

[Install]
WantedBy=multi-user.target
